# docker/barman/cron.d/barman
* * * * * barman /usr/local/bin/barman receive-wal --create-slot database; /usr/local/bin/barman cron
0 * * * * barman /usr/local/bin/barman backup all
* * * * * barman /usr/local/bin/barman receive-wal --create-slot database; /usr/local/bin/barman cron
0 * * * * barman /usr/local/bin/barman backup all
* * * * * barman /usr/local/bin/barman receive-wal --create-slot database; /usr/local/bin/barman cron
0 * * * * barman /usr/local/bin/barman backup all
* * * * * barman /usr/local/bin/barman receive-wal --create-slot database; /usr/local/bin/barman cron
0 * * * * barman /usr/local/bin/barman backup all
* * * * * barman /usr/local/bin/barman receive-wal --create-slot database; /usr/local/bin/barman cron
0 * * * * barman /usr/local/bin/barman backup all
* * * * * barman /usr/local/bin/barman receive-wal --create-slot database; /usr/local/bin/barman cron
0 * * * * barman /usr/local/bin/barman backup all
* * * * * barman /usr/local/bin/barman receive-wal --create-slot database; /usr/local/bin/barman cron
0 * * * * barman /usr/local/bin/barman backup all
* * * * * barman /usr/local/bin/barman receive-wal --create-slot database; /usr/local/bin/barman cron
0 * * * * barman /usr/local/bin/barman backup all
